# systemd/aether-discovery.service
[Unit]
Description=AETHER DMX Node Discovery Service
After=network.target

[Service]
Type=simple
User=ramzt
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/ramzt/aether-discovery.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
