# systemd/aether-core.service
[Unit]
Description=AETHER Core - Unified Control System
After=network.target

[Service]
Type=simple
User=ramzt
WorkingDirectory=/home/ramzt
ExecStart=/usr/bin/python3 /home/ramzt/aether-core.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
