# systemd/dmx-node-backend.service
[Unit]
Description=AETHER DMX Node Discovery Backend
After=network.target

[Service]
Type=simple
User=ramzt
WorkingDirectory=/home/ramzt/Aether-DMX/node-backend
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
